# .config/systemd/user/urxvtd.service
[Unit]
Description=rxvt-unicode

[Service]
ExecStart=/usr/bin/urxvtd
Restart=always

[Install]
WantedBy=default.target
